# backend/ai-backend.service
[Unit]
Description=AI Backend API Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ai/backend
ExecStart=/usr/bin/python3 /opt/ai/backend/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-backend

Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
